# build/supervisord.conf
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 3
loglevel = debug

[program:httpd]
command = /bin/bash -c "rm -rf /run/httpd/* /tmp/httpd* && exec /usr/sbin/httpd -DFOREGROUND"
stdout_logfile = /var/log/supervisor.%(program_name)s.log
stderr_logfile = /var/log/supervisor.%(program_name)s.log
autostart = true
autorestart = true

[program:crond]
command = /bin/bash -c "source /etc/sysconfig/crond && exec /usr/sbin/crond -n $CRONDARGS"
stdout_logfile = /var/log/supervisor.%(program_name)s.log
stderr_logfile = /var/log/supervisor.%(program_name)s.log
autostart = true
autorestart = true
